# certbot/certbot-renew.timer
[Unit]
Description=Renew Certbot certificates every day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
